# scripts/launchd/com.example.pa-briefing-evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- See com.example.pa-briefing-morning.plist for placeholder conventions. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.pa-briefing-evening</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_DIR__/scripts/run_briefing.sh</string>
    <string>evening</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>USER</key>
    <string>__USER__</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>19</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__REPO_DIR__/briefings/cron-evening.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_DIR__/briefings/cron-evening.err.log</string>
</dict>
</plist>
